# root/etc/uci-defaults/30_luci-theme-atmaterial
#!/bin/sh
uci batch <<-EOF
	set luci.themes.ATMaterial=/luci-static/atmaterial_pink
	commit luci
EOF
exit 0
